# deploy/systemd/rescue-mesh-api.service
# rescue-mesh-api.service: authenticated plane, HTTPS 8443 (file 01 step 8).
# TLS cert/key are issued by the fleet CA during setup (file 09 plane 2).
[Unit]
Description=Rescue mesh authenticated API (HTTPS 8443)
After=network.target rescue-mesh-runtime.service
Wants=rescue-mesh-runtime.service

[Service]
User=drone
Group=drone
WorkingDirectory=/home/drone/rescue-mesh/backend
EnvironmentFile=/etc/rescue-mesh/node.env
ExecStart=/home/drone/rescue-mesh/backend/.venv/bin/python api.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
